# every-minute.crontab
PATH=/sbin:/bin:/usr/sbin:/usr/bin

* * * * * <PATH-TO-SCRIPT>/print-connectivity.sh >> <PATH-TO-SCRIPT>/network-connectivity.log